# scripts/button-action.service
[Unit]
Description=Button action: single=Reboot, double=Wi-Fi reconnect
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/button-action.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target